# profile.d/python.sh
export PATH="/usr/local/opt/python/libexec/bin:$PATH"

export PYTHON_BASE=`python -m site --user-base`
export PATH="$PYTHON_BASE/bin:$PATH"